# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main -c config_file=/etc/postgresql/15/main/postgresql.conf -c listen_addresses=127.0.0.1
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log

[program:redis]
command=redis-server --appendonly yes --dir /var/lib/redis
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log

[program:backend]
command=python3 -m uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/backend
autostart=true
autorestart=true
priority=30
startsecs=5
startretries=10
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
environment=DATABASE_URL="postgresql://postgres@localhost:5432/ecommerce_dev",REDIS_URL="redis://localhost:6379/0",JWT_SECRET="dev-secret-key-change-in-production",ENVIRONMENT="development",CORS_ORIGINS="http://localhost:3000,http://localhost",RESEND_API_KEY="",EMAIL_FROM="onboarding@resend.dev",FRONTEND_URL="http://localhost"

[program:frontend]
command=/app/frontend/node_modules/.bin/next start -p 3000
directory=/app/frontend
autostart=true
autorestart=true
priority=30
startsecs=5
startretries=10
stdout_logfile=/var/log/supervisor/frontend.log
stderr_logfile=/var/log/supervisor/frontend_err.log
environment=NEXT_PUBLIC_API_URL="http://localhost/api/v1",NODE_ENV="production"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log
